# machines/ux32ln/rootfs/etc/sv/brightnessctl/run
#!/bin/sh
set -e
env XDG_RUNTIME_DIR=/var/lib brightnessctl -r > /dev/null
exec chpst -b brightness pause
