# rootfs/etc/cron.daily/localroute
#!/bin/sh
exec curl https://raw.githubusercontent.com/herrbischoff/country-ip-blocks/master/ipv4/ir.cidr -o /etc/localroute.d/ir.cidr -sS --fail
