# storage/etc/sv/backlight-ctl/finish
#!/bin/sh
/usr/local/bin/backlight-ctl mon save
/usr/local/bin/backlight-ctl kbd save
