# storage/etc/sv/backlight-ctl/run
#!/bin/sh
/usr/local/bin/backlight-ctl mon restore
/usr/local/bin/backlight-ctl kbd restore
exec chpst -b backlight-ctl pause
