# rootfs/etc/sv/dnsmasq/run
#!/bin/sh
mkdir -p /var/lib/misc
mkdir -p /run/resolvconf/ && touch /run/resolvconf/resolv.conf
[ -r ./conf ] && . ./conf
exec dnsmasq -k ${OPTS:---enable-dbus -u dnsmasq -g dnsmasq} 2>&1
